# roles/content_extraction/templates/gazette_content_extraction.timer.j2
[Unit]
Description=Run Gazette content extraction

[Timer]
OnCalendar={{TIMER.ONCALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
